# debian/vdr-addon-avahi-linker/etc/init/avahi-linker.conf
start on started avahi-daemon
stop on stopping avahi-daemon

respawn
kill timeout 20
kill signal INT

script

if [ -e /etc/default/locale ]; then
. /etc/default/locale
fi

if [ -n "$LANG" ]; then
   export LANG
fi

exec /usr/bin/avahi-linker -l /var/log/avahi-linker.log
end script
