# debian/prevent-umount-on-pause.service
[Unit]
After=avahi-linker.service

[Service]
ExecStart=/usr/bin/prevent-umount-on-pause

[Install]
WantedBy=multi-user.target
