# etc/init/prevent-umount-on-pause.conf
start on started avahi-linker
stop on stopping avahi-linker

script
    prevent-umount-on-pause
end script
