# etc/init/vdr-net-monitor.conf
author "Alexander Grothe <seahawk1986@gmx.de>"
description "Avahi linker"

start on started vdr
stop on stopping vdr

respawn

script
exec /usr/bin/vdr-net-monitor
end script

